# build/service/apache2/run
#!/bin/sh

set -e

_setSessionManager () {
  host=$1
  port=$2
  database=$3

  echo "host=$host"
  echo "port=$port"
  echo "database=$database"

  # PHP
  cat /etc/php/7.0/apache2/php.ini | \
    sed "s/session.save_handler = files/session.save_handler = redis\nsession.save_path = \"tcp:\/\/$host:$port\"\nextension=redis.so/" \
    > /etc/php/7.0/apache2/php.ini.tmp
  mv /etc/php/7.0/apache2/php.ini.tmp /etc/php/7.0/apache2/php.ini
}

setSessionManager () {
  if [ ! -z $REDIS_PORT_6379_TCP_ADDR ]; then
    REDIS_DATABASE=${REDIS_DATABASE:-"0"}
    _setSessionManager $REDIS_PORT_6379_TCP_ADDR $REDIS_PORT_6379_TCP_PORT ${REDIS_DATABASE:-"0"}
  elif [ ! -z $REDIS_HOST ]; then
    REDIS_DATABASE=${REDIS_DATABASE:-"0"}
    REDIS_PORT=${REDIS_PORT:-"6379"}
    _setSessionManager $REDIS_HOST $REDIS_PORT ${REDIS_DATABASE:-"0"}
  else
    echo "Warn: no redis session manager."
  fi
}

# Session manager
setSessionManager

/etc/init.d/apache2 start
